# ops/backup/schulsani-backup.service
[Unit]
Description=Verschluesselte Sicherung der SchulSani-Datenbank
After=postgresql.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/schulsani-backup.sh
